# ngrok.conf
[program:ngrok]
directory=/home/ntrllog/CS5540-HoneyGen
command=ngrok http 8000 --log=stdout
autostart=true
autorestart=false
stdout_logfile=/var/log/ngrok.log
redirect_stderr=true
