# ngrok.conf
[program:ngrok]
directory=/flaskproj
command=ngrok http 5000 --log=stdout
autostart=true
autorestart=false
stdout_logfile=/var/log/ngrok.log
redirect_stderr=true