# site2/localtunnel.conf
[program:localtunnel]
directory=/home/ntrllog/CS5540-HoneyGen/site2
command=lt -p 8000
autostart=true
autorestart=false
stdout_logfile=/var/log/localtunnel.log
redirect_stderr=true
